# supervisord.conf
[program:node]
command = node web.js
directory = /home/dotcloud/current
